# vscode_build/supervisord.conf
[supervisord]
nodaemon=true

[program:bash]
command=/bin/bash -c "while true; do sleep 1; done"
autostart=true
autorestart=true